# roles/borgbackup/files/backup.timer
[Unit]
Description=Backup dayly

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
